# bot.conf
[program:pin_bot]
directory=/root/pin_and_mailing_bot
command = python3.9 main.py

autostart=true
autorestart=true
stderr_logfile=/root/pin_and_mailing_bot/logfile.log

